# deploy/systemd/ouroboros-ibkr-critic.service
[Unit]
Description=Ouroboros IBKR Critic — weekly US stock performance evaluator
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/trading_bot/MAIN
ExecStart=/home/ubuntu/trading_bot/MAIN/.venv/bin/python tools/weekly_ibkr_critic.py
StandardOutput=append:/home/ubuntu/trading_bot/MAIN/tools/ibkr_critic.log
StandardError=append:/home/ubuntu/trading_bot/MAIN/tools/ibkr_critic.log
