# deploy/systemd/ouroboros-morning-start-check.timer
[Unit]
Description=Run Ouroboros morning safety check every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=ouroboros-morning-start-check.service

[Install]
WantedBy=timers.target
