# deploy/systemd/ouroboros-ibkr-gateway-watch.service
[Unit]
Description=Ouroboros IB Gateway watchdog for VM dashboard
After=network-online.target ouroboros-unified-dashboard.service
Wants=network-online.target
ConditionPathExists=/home/ubuntu/trading_bot/MAIN/tools/ibkr_gateway_watch.py

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/trading_bot/MAIN
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/ubuntu/trading_bot/MAIN/tools/ibkr_gateway_watch.py --url http://127.0.0.1:8793/tools/unified_dashboard.html --timeout-sec 15
StandardOutput=append:/home/ubuntu/trading_bot/MAIN/review_out/ibkr_gateway_watch.systemd.out.log
StandardError=append:/home/ubuntu/trading_bot/MAIN/review_out/ibkr_gateway_watch.systemd.err.log
