# deploy/systemd/ouroboros-weekly-autotrain.timer
[Unit]
Description=Run Ouroboros weekly report + AI auto feedback every Monday

[Timer]
OnCalendar=Mon *-*-* 00:20:00
Unit=ouroboros-weekly-autotrain.service
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
