# tools/com.ouroboros.chart_ai_predict.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ouroboros.chart_ai_predict</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/tani/trading_bot/trading_bot/MAIN/local_ai/chart_ai/predict.py</string>
    <string>--verbose</string>
  </array>

  <!-- 米国市場開始前（ET 9:00 = JST 22:00）に毎日更新 -->
  <key>StartCalendarInterval</key>
  <array>
    <!-- 平日22:00 JST (ET 9:00 開場前) -->
    <dict>
      <key>Hour</key>   <integer>22</integer>
      <key>Minute</key> <integer>0</integer>
      <key>Weekday</key><integer>1</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>22</integer>
      <key>Minute</key> <integer>0</integer>
      <key>Weekday</key><integer>2</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>22</integer>
      <key>Minute</key> <integer>0</integer>
      <key>Weekday</key><integer>3</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>22</integer>
      <key>Minute</key> <integer>0</integer>
      <key>Weekday</key><integer>4</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>22</integer>
      <key>Minute</key> <integer>0</integer>
      <key>Weekday</key><integer>5</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tani/trading_bot/trading_bot/MAIN/local_ai/chart_ai</string>

  <key>StandardOutPath</key>
  <string>/tmp/chart_ai_predict.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/chart_ai_predict.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
