# tools/com.ouroboros.ibkr_snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ouroboros.ibkr_snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tani/trading_bot/trading_bot/MAIN/tools/run_ibkr_snapshot.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/tani/trading_bot/trading_bot/MAIN/review_out/ibkr_snapshot_launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tani/trading_bot/trading_bot/MAIN/review_out/ibkr_snapshot_launchd.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
